# scripts/alpine-init/fitit-api
#!/sbin/openrc-run

name=fitit-api
cfgfile=/etc/fitit-api.env
command="NODE_ENV=production NODE_PATH=/usr/local/lib/fitit-api/node_modules DOTENV_CONFIG_PATH=/etc/fitit-api.env node /usr/local/bin/fitit-api/server.js"
pidfile="/run/fitit-api.pid"
command_background="yes"

depend() {
  need net
}