# scripts/alpine-init/fitit-api
#!/sbin/openrc-run

name=fitit-api
cfgfile=/etc/fitit-api.env
command="node NODE_ENV=production dotenv_config_path=/etc/fitit-api.env /usr/local/bin/fitit-api/server.js"
pidfile="/run/fitit-api.pid"
command_background="yes"

depend() {
  need net
}